# xe-base-config/lib/systemd/system/connection-online.target
[Unit]
Description=Common connection status online
#Wants=network-online.target
#After=network-online.target
