# xe-base-config/lib/systemd/system/connection-online@.target
[Unit]
Description=Interface %I online-status

