# xe-base-config/etc/profile.d/gdk-scroll-fix.sh
#export GDK_CORE_DEVICE_EVENTS=1
export NO_AT_BRIDGE=1

